# rootfs/etc/systemd/system/init_system.service
[Unit]
Description=DART Init-System Service
After=multi-user.target

[Service]
User=zre
TTYPath=/dev/tty1
Environment="XDG_RUNTIME_DIR=/home/zre/"
ExecStart=/usr/bin/cage -- /home/zre/init_system/bin/init-system /home/zre/init_system/bin/pre-exec /home/zre/init_system/bin/exec-can-mdf-logger /home/zre/init_system/bin/exec-dashboard-gui
StandardError=journal

[Install]
WantedBy=default.target
